# bosun/supervisord.conf
[supervisord]
nodaemon=true

[program:bosun]
command=/bosun/bosun -c /data/bosun.toml
priority=20
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0